# rc/dnsmasq/dnsmasq-net-1.in
#! /bin/sh
#
# Startup script for dnsmasq to assign static ip for libvirt domains.
#
# chkconfig init header:
#
# dnsmasq-net-1: dnsmasq init script for net-1
# chkconfig: - 99 01
# description: This script starts your DNS caching server
# processname: dnsmasq
# pidfile: /var/run/dnsmasq-net-1.pid
#
# LSB init header:
#
### BEGIN INIT INFO
# Provides: dnsmasq-net-1
# Required-Stop:
# Required-Start: libvirtd
# Default-Start: 3 4 5
# Default-Stop: 0 1 6
# Should-Start: $network
# Short-Description: dnsmasq for libvirt networks
# Description: This is a dnsmasq init script to assign ip for libvirt domains
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions


RETVAL=0

EXEC=/usr/sbin/dnsmasq
if test ! -x $EXEC; then
    echo -n $"$EXEC does not exist or is not a executable."; warning; echo
    exit 5
fi


NETWORK=net-1

prefix=@prefix@
CONFIG=@localstatedir@/lib/@PACKAGE@/networks/dnsmasq/$NETWORK.conf
HOSTSFILE=@localstatedir@/lib/@PACKAGE@/networks/dnsmasq/$NETWORK.hostsfile
PIDFILE=@localstatedir@/run/dnsmasq-$NETWORK.pid
LOCKFILE=@localstatedir@/lock/subsys/dnsmasq-net-1  # suppress rpmlint errors

OPTIONS="--conf-file=$CONFIG --pid-file=$PIDFILE --dhcp-hostsfile=$HOSTSFILE"

if test ! -f $CONFIG; then
    echo -n $"$0: CONFIG '$CONFIG' does not exist."; warning; echo
    exit 5
fi
if test ! -s $HOSTSFILE; then
    echo -n $"$0: dhcp hostsfile '$HOSTSFILE' does not exist or zero contents."; warning; echo
fi


case "$1" in
  start)
        echo -n "Starting dnsmasq for $NETWORK: "
        daemon $EXEC $OPTIONS
	RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch $LOCKFILE
        ;;
  stop)
	status -p $PIDFILE dnsmasq
        ret=$?

        if test $ret -eq 0; then
            echo -n "Shutting down dnsmasq for $NETWORK: "
            killproc -p $PIDFILE dnsmasq 
        fi
	RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && rm -f $LOCKFILE $PIDFILE
        ;;
  status)
	status -p $PIDFILE dnsmasq
	RETVAL=$?
	;;
  reload)
	echo -n "Reloading dnsmasq for $NETWORK: "
	killproc -p $PIDFILE dnsmasq -HUP
	RETVAL=$?
	echo
	;;
  restart)
	$0 stop && $0 start
	RETVAL=$?
	;;
  condrestart)
	status -p $PIDFILE dnsmasq
        ret=$?

        if test $ret -eq 0; then
            $0 stop && $0 start
            RETVAL=$?
        fi
	;;
  *)
        echo "Usage: $0 {start|stop|restart|reload|condrestart|status}"
        exit 1
esac

exit $RETVAL

